# examples/cron/workcue-weekday.cron
# Run WorkCue at 09:00 every weekday.
# Replace /path/to/workcue with your local checkout.
0 9 * * 1-5 cd /path/to/workcue && pnpm today --config .workcue/config.yml >> .workcue/workcue.log 2>&1
